# debian/python3-cyberfusion-work-item-automations.cron.d
* * * * * glwia /usr/bin/tuxis-cron /etc/cronscripts/glwia.cron
